# CTF_G3/config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/sshd.err.log
stdout_logfile=/var/log/supervisor/sshd.out.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log

[program:php-fpm]
command=/usr/sbin/php-fpm7.4 -F
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/phpfpm.err.log
stdout_logfile=/var/log/supervisor/phpfpm.out.log

[program:vsftpd]
command=/usr/sbin/vsftpd /etc/vsftpd.conf
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/vsftpd.err.log
stdout_logfile=/var/log/supervisor/vsftpd.out.log

[program:sysmon-watchdog]
command=/opt/sysmon_watchdog.sh
autostart=true
autorestart=true
user=root
stderr_logfile=/var/log/supervisor/watchdog.err.log
stdout_logfile=/var/log/supervisor/watchdog.out.log
